# www/caddy/src/etc/rc.d/caddy
#!/bin/sh
#
# PROVIDE: caddy
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# DO NOT EDIT THIS FILE -- OPNsense auto-generated file
#

. /etc/rc.subr

name="caddy"
rcvar=caddy_enable

load_rc_config $name

: ${caddy_enable:=NO}
: ${caddy_config_file:="/usr/local/etc/caddy/Caddyfile"}

pidfile="/var/run/caddy.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r /usr/local/bin/caddy run --config ${caddy_config_file} >> /var/log/caddy/caddy.log 2>&1"

# Export XDG_CONFIG_HOME
export XDG_CONFIG_HOME="/usr/local/etc/caddy/.config"

# Path to setup script
caddy_setup_script="/usr/local/opnsense/scripts/Pischem/Caddy/setup.sh"

# Define pre-command functions if needed
# caddy_stop_precmd()
# {
#     # Add any logic needed before stopping Caddy here
# }

start_precmd=caddy_start_precmd
restart_precmd=caddy_start_precmd

caddy_start_precmd()
{
    # Call the setup script before starting the service
    sh ${caddy_setup_script}
}

run_rc_command "$1"
